# Chapter 7: set-txpower.sh and set-txpower.service
#### set-txpower.sh ####
!/bin/bash
# Wait until wlan0 is up
while ! ip link show wlan0 | grep -q "state UP"; do
    sleep 1
done
# Set Tx power to 3 dBm (300 centi-dBm)
iw dev wlan0 set txpower fixed 300

#### set-txpower.service ####
[Unit]
Description=Set Wi-Fi Tx Power to 3 dBm
After=network.target
Wants=network.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/set-txpower.sh
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
